# sync/com.birdsite.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.birdsite.sync</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Replace USERNAME with your macOS username -->
    <string>/Users/USERNAME/.pyenv/versions/3.13.9/bin/python3</string>
    <string>/Users/USERNAME/repos/bird-photo-website/sync/sync.py</string>
  </array>

  <!-- Run daily at 06:00 AM -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Also trigger automatically when Photos library changes (new photo added) -->
  <key>WatchPaths</key>
  <array>
    <string>/Users/USERNAME/Pictures/Photos Library.photoslibrary/database/Photos.sqlite</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/USERNAME/Library/Logs/birdsite-sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/Library/Logs/birdsite-sync-error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/USERNAME</string>
    <key>PATH</key>
    <string>/Users/USERNAME/.pyenv/versions/3.13.9/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Throttle: don't run more than once every 5 minutes even if Photos changes repeatedly -->
  <key>ThrottleInterval</key>
  <integer>300</integer>

  <key>KeepAlive</key>
  <false/>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
